# server/deploy/rsschat.service
[Unit]
Description=rss.chat server (rssnetwork.js)
After=network.target mariadb.service
Wants=mariadb.service

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/rss.chat/server/code
ExecStart=/usr/bin/node rssnetwork.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rsschat

[Install]
WantedBy=multi-user.target
